# deploy/systemd/vapor-state-export.service
[Unit]
Description=Export Vapor server state bundle
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/vapor-server/root.env
ExecStart=/opt/vapor-server-root/deploy/scripts/export-state.sh
